# systemd/dns-hourly-export.service
# Install: /etc/systemd/system/dns-hourly-export.service
#
#   sudo install -m 0644 dns-hourly-export.service /etc/systemd/system/
#   sudo systemctl daemon-reload

[Unit]
Description=Write hourly deduplicated DNS audit (systemd-resolved journal)
Wants=systemd-resolved.service
After=systemd-resolved.service

[Service]
Type=oneshot
ExecStart=/usr/local/lib/dns-proxmox-audit/dns-hourly-export.py --output-dir /var/lib/dns-audit
# Need journal read for any host user:
#   User=root
User=root
