# systemd/dns-hourly-export.service
# Ship to /etc/systemd/system/dns-hourly-export.service

[Unit]
Description=Write hourly deduplicated DNS audit (systemd-resolved journal)
Wants=systemd-resolved.service
After=systemd-resolved.service

[Service]
Type=oneshot
Environment=PYTHONPATH=/usr/local/lib
ExecStart=/usr/bin/python3 -m dns_proxmox_audit.hourly_export --previous-hour
# Need journal read for any host user:
#   User=root
User=root
