# systemd/dns-hourly-export.timer
# Ship to /etc/systemd/system/dns-hourly-export.timer

[Unit]
Description=Run dns-hourly-export every hour (previous local hour)

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=1min
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
