# systemd/systemd-resolved.service.d/10-dns-audit.conf
# Install on the machine that runs systemd-resolved (auditing DNS).
#
# Requires: systemd 253+ for LogFilterPatterns= (keeps the journal from storing
# every debug line; tune patterns after: journalctl -u systemd-resolved -n 300).
# If the service fails to start or the key is unknown, use 10-dns-audit-no-logfilter.conf
# instead, or remove all LogFilterPatterns= lines below.

[Service]
Environment=SYSTEMD_LOG_LEVEL=debug
# Any MESSAGE matching one of these (ERE) is kept in the journal; others are not stored.
# Adjust to match lines that contain the qname and/or A/AAAA answers you care about.
LogFilterPatterns=IN A[[:space:]]
LogFilterPatterns=IN[[:space:]]AAAA[[:space:]]
LogFilterPatterns=lookup[[:space:]]key[[:space:]]
LogFilterPatterns=Looking[[:space:]]up[[:space:]]
LogFilterPatterns=Varlink[[:space:]]
LogFilterPatterns=Received[[:space:]]DNS[[:space:]]reply
LogFilterPatterns=cache[[:space:]]for[[:space:]]
LogFilterPatterns=Transaction[[:space:]]
