# systemd/systemd-resolved.service.d/10-dns-audit-no-logfilter.conf
# Install as 10-dns-audit.conf (rename) when LogFilterPatterns= is not supported
# (systemd < 253) or you want the full debug stream in the journal and will
# filter only in dns-hourly-export.py instead.
#
#   sudo install -d -m 0755 /etc/systemd/system/systemd-resolved.service.d
#   sudo install -m 0644 10-dns-audit-no-logfilter.conf /etc/systemd/system/systemd-resolved.service.d/10-dns-audit.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart systemd-resolved

[Service]
Environment=SYSTEMD_LOG_LEVEL=debug
