# systemd/systemd-resolved.service.d/10-dns-audit-no-logfilter.conf
# Use as 10-dns-audit.conf (rename) when LogFilterPatterns= is not supported
# (systemd < 253) or you want the full debug stream in the journal and will
# filter only in dns-hourly-export.py instead.

[Service]
Environment=SYSTEMD_LOG_LEVEL=debug
